# ansible/roles/frontend/files/btlgeo.conf
description "Below The Line Geo Server"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /home/btl/gunicorn-geo.sh

